# web/templates/supervisord.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn {{ flask_app }} -c {{ code_dir }}/gunicorn.conf.py
directory={{ code_dir }}
user=root
stdout_logfile=/var/log/gunicorn-access.log
stderr_logfile=/var/log/gunicorn-error.log
autostart=true
autorestart=true
environment=PATH="/usr/local/bin:/usr/bin:/bin"
